# data/cron.d/100-docker-cron
* * * * * root curl -u federation:secret -X GET http://storage.localhost/federate/cleanup >> /data/poirot/cron.log 2>&1
